# deb/DEBIAN/prerm
#!/bin/sh
set -e

systemctl --system disable rest4mqtt >/dev/null || true
systemctl --system stop rest4mqtt >/dev/null || true
